# systemd/cronicle.service
[Unit]
Description=Cronicle
After=syslog.target network.target

[Service]
Type=forking
User=root
PIDFile=/opt/cronicle/logs/cronicled.pid
ExecStart=/opt/cronicle/bin/control.sh start
ExecStop=/opt/cronicle/bin/control.sh stop

[Install]
WantedBy=multi-user.target